# systemd/mapd_server@.service.in
[Unit]
Description=MapD database server for %i
After=network.target remote-fs.target

[Service]
User=@MAPD_USER@
Group=@MAPD_GROUP@
WorkingDirectory=@MAPD_PATH@
ExecStart=@MAPD_PATH@/bin/mapd_server --config @MAPD_STORAGE@/%I/mapd.conf
KillMode=control-group
SuccessExitStatus=143
LimitNOFILE=65536
Restart=always

[Install]
WantedBy=multi-user.target
