# systemd/mapd_xorg.service
[Unit]
Description=Xorg
Before=xorg.target

[Service]
ExecStart=/usr/bin/Xorg :1 -nolisten tcp -noreset

[Install]
WantedBy=multi-user.target
