# systemd/mapd_server_render.service.in
[Unit]
Description=MapD database server
After=network.target remote-fs.target
Requires=mapd_xorg.service

[Service]
User=@MAPD_USER@
Group=@MAPD_GROUP@
WorkingDirectory=@MAPD_PATH@
Environment=DISPLAY=:1
ExecStart=@MAPD_PATH@/bin/mapd_server --config @MAPD_STORAGE@/mapd.conf
KillMode=control-group
LimitNOFILE=65536
Restart=always

[Install]
WantedBy=multi-user.target
