# systemd/mapd_server.service.in
[Unit]
Description=MapD database server
After=network.target remote-fs.target

[Service]
#User=mapd
#Group=mapd
Environment="LD_LIBRARY_PATH=/usr/local/cuda/nvvm/lib64"
WorkingDirectory=@MAPD_PATH@
ExecStart=@MAPD_PATH@/bin/mapd_server --config mapd.conf
KillMode=control-group

[Install]
WantedBy=multi-user.target
