# systemd/mapd_server_render.service.in
[Unit]
Description=MapD database server
After=network.target remote-fs.target
Requires=mapd_xorg.service

[Service]
User=@MAPD_USER@
Group=@MAPD_GROUP@
Environment=DISPLAY=:1
Environment="LD_LIBRARY_PATH=@MAPD_LIBJVM_DIR@"
ExecStart=@MAPD_PATH@/bin/mapd_server --config @MAPD_STORAGE@/mapd.conf
KillMode=control-group
LimitNOFILE=16384

[Install]
WantedBy=multi-user.target
