# systemd/mapd_sd_server@.service.in
[Unit]
Description=MapD StringDictionary Server for %i
After=network.target remote-fs.target

[Service]
User=@MAPD_USER@
Group=@MAPD_GROUP@
WorkingDirectory=@MAPD_PATH@
ExecStart=@MAPD_PATH@/bin/StringDictionaryServer --config @MAPD_STORAGE@/%I/mapd-sds.conf
LimitNOFILE=65536
Restart=always

[Install]
WantedBy=multi-user.target
