# systemd/mapd_web_server@.service.in
[Unit]
Description=MapD web server for %i
After=network.target remote-fs.target

[Service]
User=@MAPD_USER@
Group=@MAPD_GROUP@
WorkingDirectory=@MAPD_PATH@
ExecStart=@MAPD_PATH@/bin/mapd_web_server --config @MAPD_STORAGE@/%I/mapd.conf
Restart=always

[Install]
WantedBy=multi-user.target
