# systemd/omnisci_server@.service.in
[Unit]
Description=OmniSci database server for %i
After=network.target remote-fs.target

[Service]
User=@OMNISCI_USER@
Group=@OMNISCI_GROUP@
WorkingDirectory=@OMNISCI_PATH@
ExecStart=@OMNISCI_PATH@/bin/omnisci_server --config @OMNISCI_STORAGE@/%I/omnisci.conf
KillMode=control-group
SuccessExitStatus=143
LimitNOFILE=65536
Restart=always

[Install]
WantedBy=multi-user.target
